# etc/crontab
0 */12 * * * pwsh /app/webmonitor-pro/App.ps1 >> /var/log/webmonitor-pro.log 2&>1
