# etc/crontab
0 */6 * * * pwsh /apps/webmonitor-pro/App.ps1 | Out-File /var/log/webmonitor-pro.log -Append
